# inno.service
[Unit]
Description=Inno battery notification daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/inno
Restart=on-failure
RestartSec=5

# Resource limits for low CPU usage
CPUQuota=5%
MemoryMax=50M
Nice=19

[Install]
WantedBy=graphical-session.target
